# crontab
* * * * * /usr/local/bin/python /app/main.py >> /app/output/output.txt 2>&1
