# crontab
NOTION_API_TOKEN=${NOTION_API_TOKEN}
RESEND_API_KEY=${RESEND_API_KEY}
* * * * * /usr/local/bin/python /app/main.py >> /app/output/output.txt 2>&1
